# Advanced Bash Projects/Write a log parser with regular expressions/systemd/log-parser.timer
[Unit]
Description=Run Log Parser Daily
Requires=log-parser.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target