# Bonus  Cutting-Edge/Try immutable OSes like Fedora Silverblue/systemd/immutable-os-explorer.timer
[Unit]
Description=Run ImmutableOS Explorer Update Checker weekly
Documentation=https://github.com/username/immutable-os-explorer
Requires=immutable-os-explorer.service

[Timer]
OnBootSec=15min
OnUnitActiveSec=1week
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target